# nexcore-s-ui.service
[Unit]
Description=nexcore-s-ui Service
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/usr/local/nexcore-s-ui/
ExecStart=/usr/local/nexcore-s-ui/sui
Restart=on-failure
RestartSec=10s

# AUDIT.md H6 — systemd 加固。所有 sandbox 选项都在 sui 仍以 root 跑、
# 仍能 bind 低端口 / 创建 tun 设备的前提下保守选,**不破坏现有功能**:
#
#   NoNewPrivileges        进程不能 setuid 提权(零副作用)
#   ProtectSystem=full     /usr /boot /etc 只读;/usr/local 仍可写(cert/db/acme)
#   ProtectHome=true       /home /root /run/user 不可读,sui 不需要这些
#   PrivateTmp=true        独立 /tmp,免污染共享 tmp
#   LockPersonality        禁 personality(2) syscall
#   RestrictRealtime       禁实时调度
#   RestrictSUIDSGID       禁创建 setuid/setgid 文件
#
# 不加的(避免 break):
#   PrivateDevices         tun 入站要 /dev/net/tun
#   ProtectKernelTunables  redirect/tproxy 可能写 /proc/sys
#   User=non-root          要 bind 80/443 + ACME 写 cert,保留 root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
